# 70-kdump.rules.in
@if @ARCH@ s390 s390x
#
# For s390x the ELF header is created in the kdump kernel and therefore
# no kdump udev rules are required.
#
@else
#
# Kdump core headers needs to be regnerated if the CPUs or memory changes.
# For this, reload kdump.
#
# Novell Bug #389658
#

TEST=="/usr/bin/systemctl", GOTO="kdump_systemd"

@if @ARCH@ ppc ppc64 ppc64le
SUBSYSTEM=="cpu", ACTION=="online", PROGRAM="/etc/init.d/boot.kdump try-restart"
SUBSYSTEM=="cpu", ACTION=="offline", PROGRAM="/etc/init.d/boot.kdump try-restart"
@endif
SUBSYSTEM=="memory", ACTION=="add", PROGRAM="/etc/init.d/boot.kdump try-restart"
SUBSYSTEM=="memory", ACTION=="remove", PROGRAM="/etc/init.d/boot.kdump try-restart"

GOTO="kdump_end"

# Systemd limits service start rate, so if udev events are emitted too
# often, kdump will enter failed state, unless the counter is reset here.
LABEL="kdump_systemd"

@if @ARCH@ ppc ppc64 ppc64le
SUBSYSTEM=="cpu", ACTION=="online", RUN+="/usr/bin/systemctl reset-failed kdump", RUN+="/usr/bin/systemctl try-restart kdump"
SUBSYSTEM=="cpu", ACTION=="offline", RUN+="/usr/bin/systemctl reset-failed kdump", RUN+="/usr/bin/systemctl try-restart kdump"
@endif
SUBSYSTEM=="memory", ACTION=="add", RUN+="/usr/bin/systemctl reset-failed kdump", RUN+="/usr/bin/systemctl try-restart kdump"
SUBSYSTEM=="memory", ACTION=="remove", RUN+="/usr/bin/systemctl reset-failed kdump", RUN+="/usr/bin/systemctl try-restart kdump"

LABEL="kdump_end"
@endif
